# roles/backup_source/templates/systemd/backup-failure-slack-notify@.service.j2
# This file was generated with Ansible
[Unit]
Description=Send a Systemd notification to Slack reporting backup kit %i failed

[Service]
Type=oneshot
ExecStart=/backups/scripts/slack_notify.sh -s backup@%i.service
User=backupsys
Group=backupsys
Environment="OP_SERVICE_ACCOUNT_TOKEN={{ lookup('env', 'BACKUP_OP_SERVICE_TOKEN') }}"
SyslogIdentifier=backup-failure-slack-notify@%i

[Install]
WantedBy=multi-user.tar
