# roles/backup_source/templates/systemd/backup-monthly@.timer.j2
# This file was generated with Ansible
[Unit]
Description=Trigger a monthly backup using the kit %i, if it hasn't last completed within the last 7 days

[Timer]
Unit=backup@%i.service
OnCalendar=monthly
OnUnitInactiveSec=7 days
AccuracySec=6h
RandomizedDelaySec=3 days

[Install]
WantedBy=timers.target
