# roles/backup_source/templates/systemd/backup@.service.j2
# This file was generated with Ansible
[Unit]
Description=Backup service for kit %i
OnFailure=backup-failure-slack-notify@%i.service

[Service]
Type=oneshot
ExecStart=/bin/bash /backups/kits/%i/main.sh
WorkingDirectory=/backups
User=backupsys
Group=backupsys
Environment="OP_SERVICE_ACCOUNT_TOKEN={{ lookup('env', 'BACKUP_OP_SERVICE_TOKEN') }}"
SyslogIdentifier=backup@%i

[Install]
WantedBy=multi-user.target

