# roles/backup_source/templates/systemd/backup-hourly@.timer.j2
# This file was generated with Ansible
[Unit]
Description=Trigger a hourly backup using the kit %i, if it hasn't last completed within the last 30 minutes

[Timer]
Unit=backup@%i.service
OnCalendar=hourly
OnUnitInactiveSec=30 minutes
AccuracySec=1 minute
RandomizedDelaySec=30 minutes

[Install]
WantedBy=timers.target
